# aws/network/deb/postinst
#!/usr/bin/env bash
set -euo pipefail

systemctl enable --now metadata-secure.service
