# aws/network/deb/postinst.sh
#!/usr/bin/env bash
set -euo pipefail

systemctl preset --now metadata-secure.service
