# runner/service/deb/postinst
#!/usr/bin/env bash
set -euo pipefail

systemctl enable --now actions-runner.service
