# runner/ec2/deb/postinst
#!/usr/bin/env bash
set -euo pipefail

systemctl enable actions-runner-ec2
