# fare/create/deb/postinst
#!/usr/bin/env bash
set -euo pipefail

systemctl preset fare-create.service
