# fare/create/deb/postinst
#!/usr/bin/env bash
set -euo pipefail

systemctl enable fare-create.service
